# autostart/de.guido.asus-hub.desktop
[Desktop Entry]
Type=Application
Name=Asus Hub
Exec=asus-hub --hidden
Hidden=false
StartupNotify=false
X-GNOME-Autostart-enabled=true
