# systemd/user/caelestia-kitty-muter.path
[Unit]
Description=Watch for Caelestia Kitty theme changes

[Path]
PathChanged=/home/tony/.local/state/caelestia/theme/colors-kitty.conf

[Install]
WantedBy=default.target
