# systemd/user/caelestia-kitty-muter.service
[Unit]
Description=Mute Kitty Colors for Caelestia

[Service]
Type=oneshot
ExecStart=/home/tony/.config/caelestia/mute_kitty.py
