# systemd/user/mako.service.d/override.conf
# Make mako the lowest-priority, last-resort notification daemon.
# mako is D-Bus-activatable as org.freedesktop.Notifications; these conditions
# block it from grabbing the seat whenever a higher-priority handler is alive:
#   - caelestia shell (qs -c caelestia)  -> Hyprland session
#   - plasmashell                        -> KDE Plasma session
# All ExecCondition= must pass for mako to start, so it only runs when no other
# handler is present. The [q]s bracket trick stops pgrep from matching this
# ExecCondition's own sh process. Add more lines per process for other DEs.
[Service]
ExecCondition=/bin/sh -c '! pgrep -f "[q]s -c caelestia" >/dev/null'
ExecCondition=/bin/sh -c '! pgrep -x plasmashell >/dev/null'
